# linux/etc/cron.monthly/backup-apps
#!/bin/sh
apt-mark showmanual > /home/ash/.backup-apps.monthly

